# examples/ruf.service
[Unit]
Description=Ruf Service
After=network.target

[Service]
ExecStart=%h/.local/bin/ruf
Restart=always

[Install]
WantedBy=default.target
